# runit/gunicorn/run
#!/bin/sh
exec gunicorn -c /etc/gunicorn/gunicorn.conf wsgi:app
