# service/nginx/run
#!/bin/sh
exec /sbin/setuser nginx -c /etc/nginx/nginx.conf
