# etc/programs.ini
[program:nginx]
command=nginx
redirect_stderr=true
stdout_logfile=/var/log/nginx.log

[program:gunicorn]
command=gunicorn -c /srv/gunicorn.conf wsgi:app
redirect_stderr=true
stdout_logfile=/var/log/gunicorn.log
